# .cron
45 10 * * * brew update && brew cleanup && brew cask update && brew cask cleanup